# builder/deploy/adk-cockpit.service
[Unit]
Description=Forsch ADK Builder Cockpit (read-only dashboard, tailnet-only)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/.hermes/workspace/adk/builder
Environment=PYTHONPATH=/root/.hermes/workspace/adk/builder/src
Environment=FORSCH_ADK_WORKSPACE=/root/.hermes/workspace/adk
Environment=FORSCH_ADK_HOST=100.120.21.13
Environment=FORSCH_ADK_PORT=8780
ExecStart=/root/.hermes/workspace/adk/builder/.venv/bin/python -m forsch.adk_builder.app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
